# templates/default/vault.service.erb
[Unit]
Description=Vault Server
Documentation=https://www.vaultproject.io
Requires=network-online.target
After=network-online.target consul.service

[Service]
User=vault
Group=vault
EnvironmentFile=-<%= node['vault']['sysconfig_dir'] %>/vault
Restart=on-failure
ExecStart=<%= node['vault']['bin'] %>/vault server $OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
